# launchd/com.bigmac.trademinds.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bigmac.trademinds</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/bigmac/autonomous-trader/scripts/run_trademinds_forever.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/bigmac/autonomous-trader</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/bigmac/autonomous-trader/logs/launchd.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/bigmac/autonomous-trader/logs/launchd.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TRADEMINDS_DB</key>
    <string>/Users/bigmac/autonomous-trader/data/trader.db</string>
  </dict>
</dict>
</plist>
